# supervisord/backup-weekly.conf
[program:backup-hourly]
command=/usr/local/bin/mysqltos3 weekly
autostart=false
autorestart=false
stdout_events_enabled=true
stderr_events_enabled=true